# scripts/eniolabi-service-hub.service
[Unit]
Description=Eniolabi Service Hub
Documentation=https://eniolabi.com
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=forking
User=olabi
Group=olabi
WorkingDirectory=/var/www/eniolabi.com
Environment=NODE_ENV=production
Environment=PORT=3002
ExecStart=/usr/bin/pm2 start ecosystem.config.js --env production
ExecReload=/usr/bin/pm2 reload eniolabi-service-hub
ExecStop=/usr/bin/pm2 stop eniolabi-service-hub
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=eniolabi-service-hub

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/www/eniolabi.com

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
